# systemd/agent-panel-backend@.service
[Unit]
Description=Agent Panel backend (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/hevin/Developer/tools/agent-panel
Environment=NODE_ENV=production
Environment=AGENT_PANEL_SLOT=%i
ExecSearchPath=/home/hevin/.local/share/mise/installs/bun/latest/bin:/usr/local/bin:/usr/bin:/bin
# The dashboard bundle is built once by hot-deploy-agent-panel.sh before
# starting a backend, so start:backend skips the vite build for fast startup.
ExecStart=/bin/sh -lc 'case "%i" in blue) export PORT=7332 ;; green) export PORT=7333 ;; *) echo "unknown slot: %i" >&2; exit 2 ;; esac; exec /home/hevin/.local/share/mise/installs/bun/latest/bin/bun run start:backend'
Restart=on-failure
RestartSec=3
TimeoutStartSec=60
# Allow up to 5 min for graceful WebSocket drain on hot-deploy stop.
TimeoutStopSec=300
KillSignal=SIGTERM
SuccessExitStatus=143

[Install]
WantedBy=default.target
