# scripts/agent-panel.service
[Unit]
Description=Agent Panel (Rust + React)
After=network.target

[Service]
Type=simple
WorkingDirectory=__WORKDIR__
Environment=PORT=__PORT__
ExecStart=__BIN__
Restart=on-failure
RestartSec=3
StandardOutput=append:__LOG_FILE__
StandardError=append:__LOG_FILE__

[Install]
WantedBy=default.target
