# scripts/opencode-dashboard.service
[Unit]
Description=opencode-dashboard (Hono SSR)
After=network.target

[Service]
Type=simple
WorkingDirectory=__WORKDIR__
Environment=PATH=__BUN_BIN_DIR__:/usr/local/bin:/usr/bin:/bin
Environment=NODE_ENV=production
Environment=PORT=__PORT__
ExecStart=__BUN_BIN__ run start
Restart=on-failure
RestartSec=3
StandardOutput=append:__LOG_FILE__
StandardError=append:__LOG_FILE__

[Install]
WantedBy=default.target
